# voxera-os-scaffold/voxera-os/deploy/systemd/user/voxera-panel.service
[Unit]
Description=Voxera Panel (approvals + audit UI)
After=network.target

[Service]
Type=simple
WorkingDirectory=@VOXERA_PROJECT_DIR@
ExecStart=@VOXERA_PROJECT_DIR@/.venv/bin/voxera panel --host 127.0.0.1 --port 8844
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
