# voxera-os-scaffold/voxera-os/systemd/voxera-core.service
[Unit]
Description=Voxera Core (router/skills/audit)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/VoxeraOS/voxera-os-scaffold/voxera-os
ExecStart=%h/VoxeraOS/voxera-os-scaffold/voxera-os/.venv/bin/voxera daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
