# voxera-os-scaffold/voxera-os/systemd/voxera-core.service
[Unit]
Description=Voxera Core (router/skills/audit)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env voxera daemon
Restart=on-failure

[Install]
WantedBy=default.target
