# systemd/voxera-core.service
[Unit]
Description=Voxera Core (router/skills/audit)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/VoxeraOS
ExecStart=%h/VoxeraOS/.venv/bin/voxera daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
