# deploy/systemd/user/voxera-automation.timer
[Unit]
Description=Voxera Automation Runner Timer (periodic evaluation)

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
