# deploy/systemd/user/voxera-daemon.service
[Unit]
Description=Voxera Daemon (queue + missions)
After=network.target

[Service]
Type=simple
WorkingDirectory=@VOXERA_PROJECT_DIR@
ExecStart=@VOXERA_PROJECT_DIR@/.venv/bin/voxera daemon
Restart=on-failure
RestartSec=2
TimeoutStopSec=10

[Install]
WantedBy=default.target
