# deploy/systemd/user/voxera-automation.service
[Unit]
Description=Voxera Automation Runner (one-shot evaluation)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=@VOXERA_PROJECT_DIR@
ExecStart=@VOXERA_PROJECT_DIR@/.venv/bin/voxera automation run-due-once

[Install]
WantedBy=default.target
